# scripts/catboypage-deploy.service
[Unit]
Description=Auto-deploy catboypage on new commits

[Service]
Type=oneshot
ExecStart=/root/catboypage/scripts/deploy.sh
WorkingDirectory=/root/catboypage
